# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
environment=DISPLAY=":1",HOME="/home/vscode",USER="vscode"
user=vscode
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -rfbport 5901 -quiet
environment=DISPLAY=":1",HOME="/home/vscode",USER="vscode"
user=vscode
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

[program:websockify]
command=/usr/bin/websockify 0.0.0.0:6080 localhost:5901 --web=/usr/share/novnc
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/websockify.log
stderr_logfile=/var/log/websockify.err

[program:xfce4]
command=/usr/bin/dbus-launch --exit-with-session /usr/bin/startxfce4
environment=DISPLAY=":1",HOME="/home/vscode",USER="vscode"
user=vscode
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.err
